# services/10_sec.timer
[Unit]
Description = "Timer for service"

[Timer]

OnBootSec=3min

[Install]
WantedBy = default.target

