# services/script.timer
[Unit]
Description = "Timer for service"

[Timer]

OnBootSec=2min

[Install]
WantedBy = default.target

